# propagation/hrpu-propagation.service.template
[Unit]
Description=HamRadio-Pi Ultimate Propagation Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__PROPAGATION_DIR__
ExecStart=/usr/bin/python3 __PROPAGATION_DIR__/server.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
