# hrpu-propagation-module/hrpu-propagation.service
[Unit]
Description=HamRadio-Pi Ultimate Propagation Service
After=network-online.target
[Service]
Type=simple
User=pi
WorkingDirectory=/opt/hamradio-pi-ultimate/propagation
ExecStart=/usr/bin/python3 /opt/hamradio-pi-ultimate/propagation/server.py
Restart=on-failure
RestartSec=5
[Install]
WantedBy=multi-user.target
